# deploy/shipment-delay.service
[Unit]
Description=Shipment Delay Predictor (FastAPI/uvicorn)
After=network.target

[Service]
Type=simple
User=autoforward
Group=autoforward
WorkingDirectory=/opt/shipment-delay-predictor
EnvironmentFile=/opt/shipment-delay-predictor/.env
ExecStart=/opt/shipment-delay-predictor/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001 --workers 1
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
